# bankruptcy-admin-bot.service
[Unit]
Description=Bankruptcy Admin Bot
After=network.target

[Service]
Type=simple
User=root
Group=root
WorkingDirectory=/root/bankruptcy-bot
ExecStart=/root/bankruptcy-bot/venv/bin/python /root/bankruptcy-bot/admin_bot.py
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

# Переменные окружения
Environment=PYTHONPATH=/root/bankruptcy-bot
Environment=PYTHONUNBUFFERED=1

[Install]
WantedBy=multi-user.target


